# templates/etc/profile.d/solr.sh.j2
export PATH=$PATH:{{solr_parent_homedir}}/solr/bin
